# datanode/supervisor-datanode.ini
[program:datanode]
command=./hadoop-2.7.3/bin/hdfs --config /hadoop-2.7.3/etc/hadoop/ datanode
autostart=true
autorestart=true

