# datanode/supervisor-datanode.ini
[program:datanode]
command=/opt/hadoop/bin/hdfs --config /opt/hadoop/etc/hadoop/ datanode
autostart=true
autorestart=true

