# namenode/supervisor-namenode.ini
[program:namenode]
command=/opt/hadoop/bin/hdfs --config /opt/hadoop/etc/hadoop/ namenode
autostart=true
autorestart=true

[program:create-folder]
command=hdfs dfs -mkdir /proteus
priority=150
autostart=true
