# namenode/supervisor-namenode.ini
[program:namenode]
command=./opt/hadoop/bin/hdfs --config /opt/hadoop/etc/hadoop/ namenode
autostart=true
autorestart=true

