# resourcemanager/supervisor-resourcemanager.ini
[program:resourcemanager]
command=/opt/hadoop/sbin/yarn-daemon.sh start resourcemanager
autostart=true
autorestart=true

