# resourcemanager/supervisor-resourcemanager.ini
[program:resourcemanager]
command=./hadoop-2.7.3/sbin/yarn-daemon.sh start resourcemanager
autostart=true
autorestart=true

