# nodemanager/supervisor-nodemanager.ini
[program:nodemanager]
command=./hadoop-2.7.3/bin/yarn --config /hadoop-2.7.3/etc/hadoop/ nodemanager
autostart=true
autorestart=true

