# flink/supervisor-flink.ini
[program:flink]
command=/opt/flink/proteus-engine/build-target/bin/start-cluster.sh
autostart=true
priority=1






