# kafka/supervisor-kafka.ini
[program:kafka]
command=/opt/kafka/bin/kafka-server-start.sh /opt/kafka/config/server.properties
autostart=true
autorestart=true
priority=1

[program:kafka-topic]
command=/opt/kafka/bin/kafka-topics.sh --create --zookeeper zookeeper:2181 --replication-factor 1 --partitions 1 --topic proteus
autostart=true
priority=100

[program:kafka-ingest]
command=bash -c "/opt/kafka/bin/kafka-console-producer.sh --broker-list localhost:9092 --topic proteus < /opt/proteus-dataset/PROTEUS.csv; sleep 5s"
autostart=true
priority=300
