# kafka/supervisor-kafka.ini
[program:kafka]
command=/opt/kafka/bin/kafka-server-start.sh /opt/kafka/config/server.properties
autostart=true
autorestart=true
priority=1

[program:kafka-proteus]
command=sleep 20;bash -c "/opt/kafka/launcher.sh; pwd; sleep 6s"
autostart=true
priority=300






