# kafka/supervisor-kafka.ini
[program:kafka]
command=/opt/kafka/bin/kafka-server-start.sh /opt/kafka/config/server.properties
autostart=true
autorestart=true
priority=1

[program:kafka-create-topic-producer]
command=/opt/kafka/bin/kafka-topics.sh --create --zookeeper zookeeper:2181 --replication-factor 1 --partitions 1 --topic proteus
autostart=true
priority=100
